# iso/config/includes.chroot/etc/cron.d/init
@reboot sudo /bin/sh /etc/paas/script.sh

