# iso/config/includes.chroot/etc/cron.d/init
@reboot sudo /bin/sh /etc/paas/script.sh
* * * * * sudo /bin/sh /etc/paas/update.sh

